# rootfs/etc/profile.d/libva.sh
LIBVA_DRIVER_NAME=radeonsi
export LIBVA_DRIVER_NAME
